# docker/supervisor.conf
[program:gunicorn]
command = /usr/local/bin/gunicorn diamondserv.wsgi:application --bind unix:/home/app/app.sock

[program:nginx]
command = /usr/sbin/nginx

[supervisord]